# backend/deploy/sprintly.service.d/timeout.conf
# Drop-in overlay (not a full-unit replace): extend the stop timeout so a deploy
# SIGTERM lets in-flight generation finish/checkpoint instead of being SIGKILLed
# mid-build (which orphans the row and 502s an in-flight request). Must exceed the
# backend's drain deadline, which in turn exceeds the 120s vite-build timeout;
# systemd's 90s default would otherwise kill the process before the drain completes.
# A drop-in leaves the live unit's other [Service] settings untouched.
[Service]
TimeoutStopSec=150
